# debian/any-proxy/DEBIAN/prerm
#!/bin/bash
set -e

invoke-rc.d any_proxy stop
